# data/su.openrc.in
#!/sbin/openrc-run

command="@prefix@/@bindir@/su"
pidfile="/var/run/suidless.pid"
name="Suidless su "
description="Suidless su daemon"

