# data/su.openrc.in
#!/sbin/openrc-run

command="@prefix@/@bindir@/suidless-su"
command_args="--daemon"
command_background=true
pidfile="/var/run/suidless.pid"
name="Suidless su "
description="Suidless su daemon"

